# mihomo-watchdog.timer
[Unit]
Description=Run mihomo watchdog every 2 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=120

[Install]
WantedBy=timers.target
