# mihomo-watchdog.service
[Unit]
Description=mihomo watchdog timer
After=mihomo.service

[Service]
Type=oneshot
ExecStart=/data/mihomo/mihomo-watchdog.sh
